# ksurobot/install/files/ksurobot.service
[Service]
ExecStart={service_cmd}
Restart=always

[Install]
WantedBy=multi-user.target
