# ksurobot/install/files/ksurobot.service
[Service]
ExecStart={service_cmd}

[Install]
WantedBy=multi-user.target
